# roles/network/swag/files/vlan.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/internal-vlan.sh

[Install]
WantedBy=default.target